# autostart/xremap.desktop
[Desktop Entry]
Name=xremap
Type=Application
Exec=systemctl --user start xremap
Terminal=false
